# lang/python/supervisor/supervisod.conf
[program:hello]
command=python /home/liuy/test/py/supervisor/test.py 
autostart=true
stdout_logfile=/home/liuy/test/py/supervisor/test.log
